# etc/systemd/system/docker.service.d/openmediavault.conf
[Service]
EnvironmentFile=-/etc/default/docker
ExecStart=
ExecStart=/usr/bin/dockerd -H unix:///var/run/docker.sock $DOCKER_OPTS \
    $OMVDOCKER_IMAGE_PATH
